# recipes-packages/agent/files/cos-setup-initramfs.service
[Unit]
Description=cOS system initramfs setup
DefaultDependencies=no
Before=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=kairos-agent run-stage initramfs

[Install]
WantedBy=sysinit.target